# deploy/vega-backend.service
[Unit]
Description=Vega Market Edge Backend
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/vega-market-edge-main/vegaBackendServer
Environment="PATH=/home/ubuntu/vega-market-edge-main/vegaBackendServer/venv/bin"
EnvironmentFile=/home/ubuntu/vega-market-edge-main/vegaBackendServer/.env
ExecStart=/home/ubuntu/vega-market-edge-main/vegaBackendServer/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4

Restart=always

[Install]
WantedBy=multi-user.target
